# ansible-playbook/supervisord.conf
[program:inotif]
command=/usr/local/bin/inotif
user=inotif
autostart=true
autorestart=true
stdout_logfile=/var/log/inotif.log
stderr_logfile=/var/log/inotif.err